# kind/etc-systemd/system/apply-addons.service
# Runs once per boot inside the node, after kubelet brings the API up, and
# applies what a fresh (tmpfs) etcd lacks: storage, CoreDNS, cluster-admin RBAC.
# Logic is inline (no helper script). Idempotent — kubectl apply re-runs safely.
[Unit]
Description=Apply kind addons once the API server is up
After=kubelet.service
Wants=kubelet.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Generous timeout: ExecStartPre blocks until the API answers /healthz.
TimeoutStartSec=300
# super-admin.conf (O=system:masters) — works before any RBAC binding exists.
Environment=KUBECONFIG=/etc/kubernetes/super-admin.conf
ExecStartPre=/bin/bash -c 'until kubectl get --raw=/healthz >/dev/null 2>&1; do echo "[addons] waiting for apiserver..."; sleep 2; done'
ExecStart=/bin/bash -c 'echo "[addons] applying /addons" && kubectl apply -f /addons/ && echo "[addons] done"'

[Install]
WantedBy=multi-user.target
