# kind/etc-systemd/system/kubelet.service.d/20-apply-addons.conf
# Enables apply-addons.service without an enable-symlink (which can't be bind-
# mounted — docker resolves symlinks, and systemd ignores regular files in
# multi-user.target.wants/). kubelet is already enabled and starts at boot, so
# having it Want apply-addons.service pulls the oneshot in. apply-addons.service
# itself orders with After=kubelet.service.
[Unit]
Wants=apply-addons.service
