# minikube/etc-systemd-system/kubelet.service.d/10-kubeadm.conf
[Unit]
Wants=docker.socket

[Service]
ExecStart=
# byo-kubelet: dropped --bootstrap-kubeconfig=/etc/kubernetes/bootstrap-kubelet.conf — that TLS
# bootstrap flow is a kubeadm-ism for minting kubelet.conf on first boot. We ship a static
# kubelet.conf, so kubelet uses it directly and the bootstrap file never existed.
ExecStart=/var/lib/minikube/binaries/v1.35.1/kubelet --config=/var/lib/kubelet/config.yaml --hostname-override=minikube --kubeconfig=/etc/kubernetes/kubelet.conf --node-ip=192.168.49.2

[Install]
