# minikube/etc-systemd-system/kubelet.service.d/10-kubelet.conf
# Generic kubelet drop-in — baked into the `kubelet` image stage (workers).
# Node-agnostic on purpose: no --hostname-override (kubelet uses the container's hostname,
# set per-service in docker-compose) and no --node-ip (auto-detected from the default route).
# Workers join via standard kubelet TLS bootstrapping: --bootstrap-kubeconfig carries the
# cluster CA + a token minted by gen-pki.sh; the kubelet CSRs for its real client cert
# (auto-approved by addon RBAC, signed by kube-controller-manager) and writes the resulting
# kubeconfig to --kubeconfig. The control-plane image overrides ExecStart in
# 20-control-plane.conf (drop-ins apply alphabetically; last ExecStart wins).
[Unit]
Wants=docker.socket

[Service]
ExecStart=
ExecStart=/var/lib/minikube/binaries/v1.35.1/kubelet --bootstrap-kubeconfig=/etc/kubernetes/bootstrap-kubelet.conf --config=/var/lib/kubelet/config.yaml --kubeconfig=/etc/kubernetes/kubelet.conf
