# minikube/etc-systemd-system/kubelet.service.d/20-control-plane.conf
# Control-plane kubelet override — baked only into the `full` image stage, on top of the
# generic 10-kubelet.conf. No TLS bootstrapping here: the control-plane node boots from the
# static kubelet.conf + client cert generated at image build (chicken-and-egg: this kubelet
# launches the apiserver it would otherwise bootstrap against). Pins hostname + node IP to
# the identities baked into the PKI.
#
# --node-labels marks this the "system" node, AKS-style (the exact label AKS puts on its
# system nodepool). The frr-node chart schedules the ToR router Deployment onto it and
# keeps the per-node speaker DaemonSet off it (one bgpd per node — port 179). Kubelet
# self-labeling is allowed here: the restriction only covers *.kubernetes.io / *.k8s.io.
[Service]
ExecStart=
ExecStart=/var/lib/minikube/binaries/v1.35.1/kubelet --config=/var/lib/kubelet/config.yaml --hostname-override=minikube --kubeconfig=/etc/kubernetes/kubelet.conf --node-ip=192.168.49.2 --node-labels=kubernetes.azure.com/mode=system
