# minikube/etc-systemd-system/apply-addons.service
# Runs once per boot inside the node, after kubelet brings the API up, and applies what a
# fresh (tmpfs) etcd lacks: cilium, coredns, storage, cluster-admin RBAC. Logic is inline
# (no helper script), like kind. Idempotent — kubectl apply re-runs safely.
# systemd expands $VAR/$(...) in Exec= lines before bash sees them, so keep the commands
# $-free: hardcode the kubectl path (kicbase puts it under /var/lib/minikube/binaries, not on
# PATH) and use a $-free `until` loop. KUBECONFIG rides in via Environment= (process env, not
# a $-expansion). super-admin.conf (O=system:masters) works before any RBAC binding exists.
[Unit]
Description=Apply byo-kubelet addons once the apiserver is up
After=kubelet.service
Wants=kubelet.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=600
Environment=KUBECONFIG=/etc/kubernetes/super-admin.conf
# CILIUM toggle (default on; see Dockerfile ENV). The container's env lands on PID 1, so read
# it from /proc/1/environ (xargs -0 splits the NUL-separated entries; keeps this line $-free —
# systemd does not propagate container env to units). Exactly one CNI survives into the apply
# set: CILIUM=on => cilium; CILIUM=off => kindnet fallback (keeps nodes Ready — kubelet gates
# the Ready condition on a CNI conf existing in /etc/cni/net.d).
ExecStartPre=/bin/bash -c 'if xargs -0 -n1 </proc/1/environ | grep -qx CILIUM=off; then echo "[addons] CILIUM=off - using kindnet fallback"; rm -f /addons/cilium.yaml; else rm -f /addons/kindnet-fallback.yaml; fi'
ExecStartPre=/bin/bash -c 'until /var/lib/minikube/binaries/v1.35.1/kubectl get --raw=/healthz >/dev/null 2>&1; do echo "[addons] waiting for apiserver..."; sleep 2; done'
ExecStart=/bin/bash -c 'echo "[addons] applying /addons" && /var/lib/minikube/binaries/v1.35.1/kubectl apply -f /addons/ && echo "[addons] done"'
# BGP peering CRs live in their own dir with a retry loop: the cilium operator registers the
# BGPv2 CRDs some seconds after it starts, so the first applies can fail on unknown kinds.
ExecStart=/bin/bash -c 'until /var/lib/minikube/binaries/v1.35.1/kubectl apply -f /addons-bgp/ >/dev/null 2>&1; do echo "[addons] waiting for cilium BGP CRDs..."; sleep 5; done; echo "[addons] BGP peering applied"'

[Install]
WantedBy=multi-user.target
