# minikube/etc-systemd-system/node-routes.service
[Unit]
Description=Static routes to the FRR router (pod CIDR + foreign node subnets)
# Runs on every node (baked into the kubelet image stage). Routes land before kubelet
# workloads need them; if FRR isn't up yet that's fine — routes are just next-hops, and
# kubelet/cilium retry their connections anyway.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/local/bin/node-routes.sh

[Install]
WantedBy=multi-user.target
